# contrib/docker/prod/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server
priority=12
user=redis
group=redis

[program:tipboard]
command=/opt/tipboard/bin/tipboard runserver 0.0.0.0 7272
priority=14
